# scripts/sprintflow-celery.service
[Unit]
Description=SprintFlow Celery worker (SQS broker)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/sprintflow/backend
ExecStart=/home/ubuntu/sprintflow/venv/bin/celery -A app.worker worker --loglevel=info --concurrency=2
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sprintflow-celery

Environment=AWS_REGION=ap-south-1

[Install]
WantedBy=multi-user.target
